# scripts/com.fifthdragon.sync-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fifthdragon.sync-weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/davidliao/git_repos/fifth_dragon_capital/scripts/sync_weekly.sh</string>
    </array>

    <!-- Run every Sunday at 7:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/davidliao/git_repos/fifth_dragon_capital/logs/launchd_weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/davidliao/git_repos/fifth_dragon_capital/logs/launchd_weekly_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/davidliao/git_repos/py312/venv/bin</string>
    </dict>
</dict>
</plist>
